# supervisor/aws-name-server.conf
[program:aws-name-server]
command=aws-name-server --domain __YOUR_DOMAIN_HERE__
stderr_logfile=/var/log/aws-name-server.log
